# launchd/com.pgpclaw.openbao.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  PGPClaw — OpenBao Auto-Start + Unseal
  Install: cp launchd/com.pgpclaw.openbao.plist ~/Library/LaunchAgents/
  Load:    launchctl load ~/Library/LaunchAgents/com.pgpclaw.openbao.plist

  NOTE: Replace __REPO_DIR__ with the actual repo path during setup.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pgpclaw.openbao</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>
            docker start pgpclaw-openbao 2>/dev/null || true;
            sleep 5;
            __REPO_DIR__/openbao/scripts/unseal-bao.sh
        </string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/pgpclaw-openbao.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/pgpclaw-openbao.err</string>
</dict>
</plist>
